# systemd/workstation-workflow-backup.timer
[Unit]
Description=Run WORKSTATION1 workflow backup every 15 minutes

[Timer]
OnBootSec=3min
OnCalendar=*:0/15
AccuracySec=1min
RandomizedDelaySec=30s
Persistent=true
Unit=workstation-workflow-backup.service

[Install]
WantedBy=timers.target
